# mysql/cron/cron.d
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=""

* * * * * /usr/bin/date >>/tmp/crontest.txts
